# autorun/99-mmcautorun.rules
ACTION=="add", SUBSYSTEMS=="block", RUN+="/usr/share/scripts/mmcautorun.sh"
