# autorun/99-usbautorun.rules
ACTION=="add",KERNEL=="sd*", SUBSYSTEMS=="usb", ATTRS{product}=="Mass Storage", RUN+="/usr/share/scripts/usbautorun.sh"
